# supervisor.conf.d/03/eosio-agent.yanmei.yonghu111.conf
[program:eosio-agent.yanmei.yonghu111]
directory=/data/eosio.nocoin/agent.yanmei.yonghu111/
command=eosio-agent -d /data/eosio.nocoin/agent.yanmei.yonghu111/ --http-validate-host=false --verbose-http-errors --actor='{"name":"yonghu111","wif_keys":["5KieNtvURgF7JWWyAi2iqSYgetqNNGRbhihzZpYPRjAHjQb5A6v"]}' --nodeos-address=http://172.16.30.6/producers.http/ --relay-address=http://123.56.179.72:8080 --trace-push-address=http://123.56.179.72:8081 --http-server-address=0.0.0.0:42016
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/data/eosio.nocoin/agent.yanmei.yonghu111/run.log
environment=
