# supervisor.conf.d/03/eosio-agent.42000.conf
[program:eosio-agent.42000]
directory=/data/eosio.nocoin/agents/
command=eosio-agent -d /data/eosio.nocoin/agents/ --http-validate-host=false --verbose-http-errors --actor='{"name":"agent.main","wif_keys":["5JjRntmbxfM8zKPeAohje7ur2R5dHsbpvgUDYwbef5FNgitmRYc"]}' --multikey='["5JjRntmbxfM8zKPeAohje7ur2R5dHsbpvgUDYwbef5FNgitmRYc"]' --nodeos-address=http://172.16.30.6/producers.http/ --relay-address=http://123.56.179.72:8080 --trace-push-address=http://123.56.179.72:8081 --http-server-address=0.0.0.0:42000
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/data/eosio.nocoin/agents/42000.log
environment=
