# supervisor.conf.d/03/nodeos.backup.conf
[program:nodeos.backup]
directory=/data/eosio.nocoin/backup/
command=nodeos --plugin eosio::chain_api_plugin --plugin eosio::http_plugin --plugin eosio::net_api_plugin --chain-state-db-size-mb=2048 --access-control-allow-origin=* --contracts-console --http-validate-host=false --verbose-http-errors --max-transaction-time 900000000 --p2p-max-nodes-per-host=1000 --genesis-json /root/eosio.deploy/genesis.json -d /data/eosio.nocoin/backup/ --http-server-address=0.0.0.0:18888 --p2p-listen-endpoint=0.0.0.0:19000 --plugin eosio::history_log_plugin --history-log-dir=/data/eosio.nocoin/backup/history-log --ws-cache-endpoint=0.0.0.0:41201 --p2p-peer-address=172.16.30.6:9001 --p2p-peer-address=172.16.30.6:9002 --p2p-peer-address=172.16.30.6:9003
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/data/eosio.nocoin/backup/run.log
environment=
