# supervisor.conf.d/02/nodeos.producer111g.conf
[program:nodeos.producer111g]
directory=/data/eosio.nocoin/producers/producer111g
command=nodeos --plugin eosio::chain_api_plugin --plugin eosio::http_plugin --plugin eosio::net_api_plugin --plugin eosio::producer_api_plugin  --max-clients=0 --chain-state-db-size-mb=2048 --access-control-allow-origin=* --contracts-console --http-validate-host=false --verbose-http-errors --max-transaction-time 900000000 --p2p-max-nodes-per-host=10000 --genesis-json /root/eosio.deploy/genesis.json -e -d /data/eosio.nocoin/producers/producer111g/ --http-server-address=0.0.0.0:8007 --p2p-listen-endpoint=0.0.0.0:9007 --p2p-peer-address=127.0.0.1:9000 --p2p-peer-address=172.16.30.6:9001 --p2p-peer-address=172.16.30.6:9002 --p2p-peer-address=172.16.30.6:9003 -p producer111g --private-key [\"LCN7mkqf9YLRGGi8zM3TkkZxxMqMi7CeEJToEpehprxJj2JoBkvzR\",\"5K1ppMDm2q38oJuicQD2AjFyfJT7S9y8h1scPNKSLAYm86bNdZy\"]
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/data/eosio.nocoin/producers/producer111g/run.log
environment=
