# supervisor.conf.d/02/nodeos.producer111h.conf
[program:nodeos.producer111h]
directory=/data/eosio.nocoin/producers/producer111h
command=nodeos --plugin eosio::chain_api_plugin --plugin eosio::http_plugin --plugin eosio::net_api_plugin --plugin eosio::producer_api_plugin  --max-clients=0 --chain-state-db-size-mb=2048 --access-control-allow-origin=* --contracts-console --http-validate-host=false --verbose-http-errors --max-transaction-time 900000000 --p2p-max-nodes-per-host=10000 --genesis-json /root/eosio.deploy/genesis.json -e -d /data/eosio.nocoin/producers/producer111h/ --http-server-address=0.0.0.0:8008 --p2p-listen-endpoint=0.0.0.0:9008 --p2p-peer-address=127.0.0.1:9000 --p2p-peer-address=172.16.30.6:9001 --p2p-peer-address=172.16.30.6:9002 --p2p-peer-address=172.16.30.6:9003 -p producer111h --private-key [\"LCN8Mvva8am4acWbfWq5fFL2pArtNSe3A4oubK3DEZ6n1ogdBWMmT\",\"5KFrYvZqLn8Jqn5RzaGu37vY9uwyKzEFXYYSD2mNxBGJzVopVLp\"]
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/data/eosio.nocoin/producers/producer111h/run.log
environment=
